# build/debian/lib/systemd/system/cabby.service
[Unit]
Description=Cabby TAXII 2 Server
ConditionPathExists=/usr/bin/cabby

[Service]
ExecStart=/usr/bin/cabby
Group=cabby
Restart=on-failure
RestartSec=5
SyslogIdentifier=cabby
Type=simple
User=cabby

[Install]
WantedBy=multi-user.target
