# build/debian/etc/systemd/system/cabby.service.d/override.conf
[Service]
Environment="CABBY_ENVIRONMENT=production"
