# build/lib/systemd/system/cabby.service
[Unit]
Description=Cabby TAXII 2 Server
ConditionPathExists=/usr/local/bin/cabby
After=network.target
StartLimitInterval=10

[Service]
Type=simple
User=cabby
Group=cabby
Restart=on-failure
RestartSec=5
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/cabby
